# deploy/systemd/call-analyzer-service.service
[Unit]
Description=Call Analyzer Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=callanalyzer
Group=callanalyzer
WorkingDirectory=/opt/call-analyzer
Environment="PATH=/opt/call-analyzer/venv/bin"
EnvironmentFile=/opt/call-analyzer/.env
ExecStart=/opt/call-analyzer/venv/bin/python /opt/call-analyzer/call_analyzer/main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/call-analyzer/service.log
StandardError=append:/var/log/call-analyzer/service.error.log

[Install]
WantedBy=multi-user.target

